# systemd/post/fluidsynth/systemd/fs_mididings.service
[Unit]
Description=Mididings PATCH
Requires=sound.target jackd.service
After=sound.target jackd.service
#After=sound.target

[Service]
Type=simple
User=root
#ExecStart=/root/organ/pmi_baroque_organ-moerdijk.py
#ExecStart=/root/organ/mididings_load.py
ExecStart=/root/fluidsynth/mididings_patch.py

[Install]
WantedBy=multi-user.target
