# systemd/post/fluidsynth/systemd/fs_m-audio.path
[Unit]
Description=m-audio sound card

[Path]
DirectoryNotEmpty=/sys/class/sound/card0

#[Install]
#WantedBy=default.target
